# files/etc/uci-defaults/50-qmi-advanced
rm -f /lib/netifd/proto/qmi.sh
ln -s /lib/netifd/proto/qmi-advanced.sh /lib/netifd/proto/qmi.sh
